# transcriptr.service
[Unit]
Description=Audio Transcription Service
After=network.target

[Service]
ExecStart=/home/pi/.cache/pypoetry/virtualenvs/transcriptr-2b3aKC51-py3.11/bin/python -m transcriptr.record -d /mnt/ssd-4tb/audio-recordings/ -t 180
WorkingDirectory=/usr/local/bin
SyslogIdentifier=transcriptr
Restart=always
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
